# build/supervisord.conf
[supervisord]
nodaemon=true

[program:validator]
command=/usr/local/bin/validator --accept-terms-of-use --prater --wallet-dir /rocketpool/data/validators/prysm-non-hd --wallet-password-file /rocketpool/data/validators/prysm-non-hd/direct/accounts/secret --beacon-rpc-provider "localhost:5053" --graffiti "RP-Avado"

autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:rocketpool-node]
command=/usr/local/bin/rocketpoold --config /srv/rocketpool/config.yml --settings /srv/rocketpool/settings.yml node

autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0